# .config/autostart/terminalmax.desktop
[Desktop Entry]
Name=terminalmax
Exec=/home/dwa/bin/terminalmax
Icon=utilities-terminal
Type=Application
StartupNotify=false
